# remarkable-backup.service.template
[Unit]
Description=reMarkable tablet backup

[Service]
Type=oneshot
StandardOutput=journal
StandardError=journal
